# systemd.service
[Unit]
Description=Relay server service.

[Service]
Type=simple
ExecStart=$APP_NAME

[Install]
WantedBy=multi-user.target