# systemd.service.template
[Unit]
Description=Relay server service.

[Service]
Type=simple
ExecStart=$APP_NAME -k=${CRYPTO_KEY}

[Install]
WantedBy=multi-user.target